# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:api]
command=python -m uvicorn api:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/api.err.log
stdout_logfile=/var/log/supervisor/api.out.log
environment=PYTHONUNBUFFERED=1

[program:fetcher]
command=python fetcher.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/fetcher.err.log
stdout_logfile=/var/log/supervisor/fetcher.out.log
environment=PYTHONUNBUFFERED=1

[program:worker]
command=python worker.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/worker.err.log
stdout_logfile=/var/log/supervisor/worker.out.log
environment=PYTHONUNBUFFERED=1
